# crontab
* * * * * find /media/ -name '*.*' -mtime +90 -exec rm {} \;


